# ansible/templates/coffeepot.service.j2
[Unit]
Description=I am a CoffeePot

[Service]
ExecStart=/bin/coffeepot -distribution {{ ansible_facts['distribution'] }}
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
